# scripts/upstart.conf
# Strato - strato job file
 
description "Strato data portal"
author "WellDone International <welldone.org>"
 
# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn
 
# When to start the service
start on runlevel [2345]
 
# When to stop the service
stop on runlevel [016]
 
expect daemon;

env DATABASE_URL=mongodb://localhost:27017/strato

# Start the process
script
  cd /welldone
	exec npm start
end script